# docker/root-fs/etc/s6-overlay/s6-rc.d/setup/run
#!/command/with-contenv bash
# shellcheck shell=bash

mkdir -p /config

chown -R dload:dload /app
chown -R dload:dload /config
